# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

if [ -f /usr/bin/apt ]; then
  # ubuntu
  /usr/sbin/cron -f -L 5
fi
if [ -f /sbin/apk ]; then
  # alpine
  exec /usr/sbin/crond -f -S -l 5 -c /etc/crontabs
fi
